# package/app/luci-app-wifidogx/root/etc/uci-defaults/90_luci_wifidog
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@wifidogx[-1]
	add ucitrack wifidogx
	set ucitrack.@wifidogx[-1].init=wifidogx
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
