# package/base-files/files/etc/rc.local
# 将您的自定义命令放在此处，应该执行一次
# 系统初始化完成。默认情况下，此文件不执行任何操作。

exit 0
